# scripts/radio.service
[Unit]
Description=Zenith Radio
After=network.target pipewire-pulse.socket

[Service]
Type=simple
WorkingDirectory=/opt/radio
EnvironmentFile=/opt/radio/.env
ExecStart=/opt/radio/radio
Restart=on-failure
RestartSec=5

# Capture stdout/stderr via journald; view with: journalctl --user -u radio -f
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
