# scripts/radio-cache-cleanup.timer
[Unit]
Description=Daily Zenith Radio librespot cache cleanup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
